# templates/kibana.service
[Unit]
Description=Kibana Server

[Service]
TimeoutStartSec=0
Restart=always
ExecStart=/usr/local/bin/kibana
[Install]
WantedBy=multi-user.target
